# Actividad4/mi-servicio.service
[Unit]
Description=Servicio Actividad 4

[Service]
ExecStart=/usr/local/bin/Servicio.sh

[Install]
WantedBy=multi-user.target
